# debian/dyatel.cron.d

1 2 * * * dyatel if [ -x /usr/bin/dyatel-backup ]; then nice -n10 /usr/bin/dyatel-backup; fi

